# openrc/certmgr
#!/sbin/openrc-run

name="$SVCNAME"
supervisor="supervise-daemon"
command="/usr/local/bin/$SVCNAME"


depend() {
        need net localmount
        use logger
        after firewall
}
