# openrc/certmgr
#!/sbin/openrc-run

name="$SVCNAME"
supervisor="supervise-daemon"
command="/usr/bin/$SVCNAME"
command_args="${CERTMGR_OPTS} serve"


depend() {
        need net localmount
        use logger
        after firewall
}
